# content/service/1/run
#!/bin/bash

exec 2>&1
mkdir /data 2>/dev/null

PARSE_DB_URL() {
    # extract the protocol
    proto="$(echo $DATABASE_URL | grep '://' | sed -e's,^\(.*://\).*,\1,g')"
    if [[ "${proto}" =~ postgres ]]; then
        export DB_TYPE=postgres
        export DB_SSL_MODE=require
    elif [[ "${proto}" =~ mysql ]]; then
        export DB_TYPE=mysql
        export DB_SSL_MODE=true
    fi

    # remove the protocol
    url=$(echo $DATABASE_URL | sed -e s,${proto},,g)

    # extract the user and password (if any)
    userpass="$(echo $url | grep @ | cut -d@ -f1)"
    export DB_PASS=$(echo $userpass | grep : | cut -d: -f2)
    if [ -n "$DB_PASS" ]; then
        export DB_USER=$(echo $userpass | grep : | cut -d: -f1)
    else
        export DB_USER=$userpass
    fi

    # extract the host -- updated
    hostport=$(echo $url | sed -e s,$userpass@,,g | cut -d/ -f1)
    export DB_PORT=$(echo $hostport | grep : | cut -d: -f2)
    if [ -n "$DB_PORT" ]; then
        export DB_HOST=`echo $hostport | grep : | cut -d: -f1`
    else
        export DB_HOST=$hostport
    fi
    if [[ ${DB_TYPE} = postgres ]] && [[ ${DB_PORT} = "" ]]; then
        export DB_PORT=5432
    fi

    # extract the name (if any)
    export DB_NAME="`echo $url | grep / | cut -d/ -f2- | sed 's|?.*||'`"
}

if [ "${DATABASE_URL}" != "" ]; then
    PARSE_DB_URL
fi

EXEC=$(head /dev/urandom | md5sum | cut -c 1-9)
ln -sf /usr/bin/alist /usr/bin/${EXEC}
export HTTP_PORT=${PORT}
export TEMP_DIR=/temp
cd /
# Optional: reset admin password before starting server
# Usage: set env RESET_PASSWORD to the new password
if [ -n "${RESET_PASSWORD}" ]; then
  echo "[INFO] RESET_PASSWORD is set; trying to reset via alist CLI"
  # Ensure /data exists (you already do mkdir /data)
  # Use --data /data so CLI can find config that points to PostgreSQL
  /usr/bin/alist admin set "${RESET_PASSWORD}"  \
    && echo "[INFO] Password reset via alist CLI OK" \
    || echo "[WARN] alist CLI reset failed (missing config under /data or DB not reachable)"
fi
echo "[INFO] Checking AList database backend"

if [ -n "$DATABASE_URL" ]; then
  if [[ "$DATABASE_URL" =~ ^postgres ]]; then
    echo "[INFO] Database backend: PostgreSQL (remote)"
  elif [[ "$DATABASE_URL" =~ ^mysql ]]; then
    echo "[INFO] Database backend: MySQL (remote)"
  else
    echo "[WARN] DATABASE_URL is set but protocol is unknown"
  fi
else
  if [ -f "/data/config.json" ]; then
    DB_TYPE=$(jq -r '.database.type // "sqlite"' /data/config.json)
    echo "[INFO] Database backend (from config): $DB_TYPE"
  else
    echo "[INFO] Database backend: sqlite (default/local)"
  fi
fi
if [ -f /data/config.json ]; then
  echo "========== AList config.json =========="
  cat /data/config.json
  echo "======================================="
else
  echo "[WARN] /data/config.json not found"
fi

exec ${EXEC} server --no-prefix
# Optional: reset admin password before starting server
# Usage: set env RESET_PASSWORD to the new password
if [ -n "${RESET_PASSWORD}" ]; then
  echo "[INFO] RESET_PASSWORD is set; trying to reset via alist CLI"
  # Ensure /data exists (you already do mkdir /data)
  # Use --data /data so CLI can find config that points to PostgreSQL
  /usr/bin/alist admin set "${RESET_PASSWORD}" \
    && echo "[INFO] Password reset via alist CLI OK${RESET_PASSWORD}" \
    || echo "[WARN] alist CLI reset failed (missing config under /data or DB not reachable)"
fi

